# jupyter/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:tensorboard]
command=tensorboard --logdir=/data/tf-data

[program:jupyter]
command=jupyter notebook --ip 0.0.0.0 --allow-root
